# roles/configure_autoreboot/templates/2d-autoreboot.timer.j2
#This file was configured via ansible
[Unit]
Description= 2degrees Reboot Scheduling.

[Timer]
OnCalendar={{ vm_autoreboot_time }}
Unit=reboot.target

[Install]
WantedBy=timers.target
