# roles/configure_autoreboot/files/2d-autoreboot.timer
[Unit]
Description=Reboot Scheduling.

[Timer]
OnCalendar=REPLACEME
Unit=reboot.target

[Install]
WantedBy=timers.target
